# modules/chezfiles/files/system/chezfiles-apply.path
[Unit]
Description=Automatically apply chezfiles when they change

[Path]
PathModified=/usr/share/bluebuild/chezfiles/system/

[Install]
WantedBy=default.target