# modules/chezfiles/files/user/chezfiles-apply.path
[Unit]
Description=Automatically apply chezfiles when they change

[Path]
PathModified=/usr/share/bluebuild/chezfiles/user/

[Install]
WantedBy=default.target