# files/systemd/system/idwt-apply.path
[Unit]
Description=Run IDWT service when config changes

[Path]
PathModified=/etc/idwt/config.yml

[Install]
WantedBy=multi-user.target