# files/systemd/system/idwt-apply.service
[Unit]
Description=IDWT service

[Service]
Type=oneshot
ExecStart=/usr/bin/idwt apply

[Install]
WantedBy=default.target