# files/system/gnome/usr/etc/profile.d/ssh_auth_gcr.sh
export SSH_AUTH_SOCK=$XDG_RUNTIME_DIR/gcr/ssh
export SSH_ASKPASS=/usr/libexec/gcr-ssh-askpass
