# files/image-pinned-etcs/common/nushell/profile.d/1-paths.nu
use std *

path add '/home/linuxbrew/.linuxbrew/bin'
path add $"($env.HOME)/.local/bin/distrobox-exported/"
path add $"($env.HOME)/.local/share/flatpak/exports/bin/"
path add $"($env.HOME)/go/bin/"