# files/image-pinned-etcs/common/nushell/profile.d/999-autorun.nu
$env.config = {
  show_banner: false,
}

smallfetch --logo none

alias l = ls -l
alias ll = ls -al

use ~/.cache/starship/init.nu