# modules/homefiles/bluebuild-homefiles-apply.path
[Unit]
Description=Automatically apply chezfiles when they change

[Path]
PathModified=/usr/share/bluebuild/chezfiles/

[Install]
WantedBy=default.target